# docker/supervisor/laravel-scheduler.ini
[program:laravel-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "php /var/www/html/artisan schedule:work"
autostart=true
autorestart=true
user=inovedados
redirect_stderr=true
stdout_logfile=/var/log/inovedados/laravel-scheduler.log
stopwaitsecs=3600