# docker/supervisor/laravel-queue.ini
[program:laravel-queue]
process_name=%(program_name)s
command=/bin/bash -c "php /var/www/html/artisan queue:work"
autostart=true
autorestart=true
user=inovedados
redirect_stderr=true
stdout_logfile=/var/log/inovedados/laravel-queue.log
stopwaitsecs=3600